# timers/aws.service
[Unit]
Description=AWS login service

[Service]
Type=exec
ExecStart=aws-azure-login --profile staging --no-sandbox --no-prompt
