# timers/aws_prod.service
[Unit]
Description=AWS Prod login service

[Service]
Type=exec
ExecStart=aws-azure-login --profile production --no-sandbox --no-prompt
